# openwrt/files/etc/profile.d/advanced_banner.sh
#!/bin/sh

if [ -n "$PS1" ]; then
  /usr/bin/advanced_banner
fi