# openwrt/files/etc/uci-defaults/99-dae
uci batch << EOF
del network.globals.ula_prefix
commit
EOF